# launchd/com.harmsen.mailprocessor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.harmsen.mailprocessor</string>

  <key>WorkingDirectory</key>
  <string>/Users/hp/proj/research/mailprocessor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/hp/.local/bin/uv</string>
    <string>run</string>
    <string>main.py</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/hp</string>
    <key>PATH</key>
    <string>/Users/hp/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/hp/Library/Logs/mailprocessor/stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/hp/Library/Logs/mailprocessor/stderr.log</string>
</dict>
</plist>
